# init/retain-sync-age@.service
[Unit]
Description=Adjusts file priority for time

[Service]
Type=oneshot
ExecStart="$prefix/bin/retain-sync" age %i
