# docs/unit/retain-sync@.service.in
[Unit]
Description=Watch for file access in a specific directory
Documentation=man:retain-sync(1)
Wants=retain-sync-resync@%i.timer
Requires=retain-sync-age@%i.timer

[Service]
Type=simple
ExecStart=@DAEMON@ %i
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
