# docs/unit/zielen@.service.in
[Unit]
Description=Watch for file access in a specific directory
Documentation=man:zielen(1)

[Service]
Type=simple
ExecStart=@bindir@/zielend %i
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
