# unit/retain-sync-resync@.service
[Unit]
Description=Pulls commonly accessed files from a server

[Service]
Type=oneshot
ExecStart=/usr/local/bin/retain-sync sync %i
