# unit/retain-sync-age@.service
[Unit]
Description=Adjusts file priority for time

[Service]
Type=oneshot
ExecStart=/usr/local/bin/retain-sync age %i
